# uwebfilter.d/uwebfilterd.service
[Unit]
Description=uwebfilterd
After=network.target

[Service]
ExecStart=/usr/bin/uwebfilterd
Restart=always

[Install]
WantedBy=multi-user.target
