# uwebfilter.d/debian/DEBIAN/postinst
#!/bin/sh
set -e

systemctl enable uwebfilterd.service
systemctl start uwebfilterd.service

exit 0
